# dist/jfy_collector.conf
# jfy_collector - Push JFY Solar Inverter stats to PVOutput

description     "JFY Collector"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec /usr/local/bin/jfy_collector
